# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:minecraft]
command=java -server -Xmx1g -Xms256m -jar minecraft_server.1.6.4.jar nogui
autorestart=true
