# sv/minecraft/run
#!/bin/sh

#sv start socklog-unix || exit 1

cd /
rm -f command-fifo
mkfifo command-fifo

exec 2>&1
exec java -server -Xmx1g -Xms256m -jar /craftbukkit-*.jar nogui <> command-fifo
