# sv/frontail/run
#!/bin/bash

sv start minecraft || exit 1

source /etc/envvars

exec 2>&1
exec /node/bin/frontail --number 100 /var/log/minecraft/current
